# roles/motd/templates/update-motd.d/50-updates.j2
#!/bin/bash

num=`apt list --upgradeable | egrep '[a-zA-Z0-9]*\/[a-zA-Z0-9]*.*' | wc -l`
echo "-----------------------------"
echo "There are $num available updates."
